# files/system/etc/profile.d/neofetch.sh
#!/usr/bin/env bash

alias neofetch=fastfetch