# files/system/etc/profile.d/open.sh
#!/usr/bin/env bash

alias open="xdg-open &>/dev/null"